# log_console/supervisor_console.ini
;/etc/supervisord/console.ini
[program:console]
;bash -c 'source /home/liyan/console/env/bin/activate; python /home/liyan/console/web_console.py'
command=/srv/www/console/env/bin/python /srv/www/console/web_console.py
directory=/srv/www/console
environment=PATH="/srv/www/console/env/bin:/usr/bin"
user=liyan
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=2
stderr_logfile=/var/log/supervisor/%(program_name)s_errors.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=2
